# root/etc/services.d/sickchill/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc python /app/sickchill/SickBeard.py \
	--datadir /config
